# deploy/systemd/koffiebon-reverb.service
# Koffiebon — Reverb websocket-server
#   sudo cp deploy/systemd/koffiebon-reverb.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now koffiebon-reverb
[Unit]
Description=Koffiebon Reverb websocket server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/koffiebon
ExecStart=/usr/bin/php /var/www/koffiebon/artisan reverb:start --host=127.0.0.1 --port=8081
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
